# canyon-alerter.service
[Unit]
Description=Runs canyon-alerter
Wants=canyon-alerter.timer

[Service]
ExecStart=/tmp/canyon-alerter.sh
WorkingDirectory=/tmp
